# contrib/init/bitcoinreloadd.conf
description "BitcoinReload Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env BITCOINRELOADD_BIN="/usr/bin/bitcoinreloadd"
env BITCOINRELOADD_USER="bitcoinreload"
env BITCOINRELOADD_GROUP="bitcoinreload"
env BITCOINRELOADD_PIDDIR="/var/run/bitcoinreloadd"
# upstart can't handle variables constructed with other variables
env BITCOINRELOADD_PIDFILE="/var/run/bitcoinreloadd/bitcoinreloadd.pid"
env BITCOINRELOADD_CONFIGFILE="/etc/bitcoinreload/bitcoinreload.conf"
env BITCOINRELOADD_DATADIR="/var/lib/bitcoinreloadd"

expect fork

respawn
respawn limit 5 120
kill timeout 600

pre-start script
    # this will catch non-existent config files
    # bitcoinreloadd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$BITCOINRELOADD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run bitcoinreloadd."
        echo "The setting must appear in $BITCOINRELOADD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"BitcoinReload Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$BITCOINRELOADD_PIDDIR"
    chmod 0755 "$BITCOINRELOADD_PIDDIR"
    chown $BITCOINRELOADD_USER:$BITCOINRELOADD_GROUP "$BITCOINRELOADD_PIDDIR"
    chown $BITCOINRELOADD_USER:$BITCOINRELOADD_GROUP "$BITCOINRELOADD_CONFIGFILE"
    chmod 0660 "$BITCOINRELOADD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$BITCOINRELOADD_PIDFILE" \
    --chuid $BITCOINRELOADD_USER:$BITCOINRELOADD_GROUP \
    --exec "$BITCOINRELOADD_BIN" \
    -- \
    -pid="$BITCOINRELOADD_PIDFILE" \
    -conf="$BITCOINRELOADD_CONFIGFILE" \
    -datadir="$BITCOINRELOADD_DATADIR" \
    -disablewallet \
    -daemon

